# systemd/storage-breaker-logrotate.service
[Unit]
Description=Rotate Storage Breaker application logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/storage-breaker
